# adscrawler/dramatiq/crawl-store-apps-apple@.service
# systemd template unit for an Apple App Store Dramatiq worker.
#
# Instantiate per country-priority group.  The group number (1 or 2) is
# passed via %i.
#
# Examples:
#
#   sudo systemctl enable --now crawl-store-apps-apple@1
#   sudo systemctl enable --now crawl-store-apps-apple@2
#
# Apple scraping is network-heavy (slower API responses), so more threads
# per process (-t 4) allow concurrent HTTP requests.  Fewer processes (-p 2)
# keep the overall concurrency reasonable.
#
# The %i placeholder is expanded to the instance name (e.g. "1").

[Unit]
Description=Store crawler – Apple group %i
After=network.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1

Environment=dramatiq_prom_port=9192
Environment=dramatiq_prom_db=/tmp/dramatiq-prometheus-apple-%i
Environment=prometheus_multiproc_dir=/tmp/dramatiq-prometheus-apple-%i
ExecStartPre=/bin/sh -c 'mkdir -p /tmp/dramatiq-prometheus-apple-%i && rm -rf /tmp/dramatiq-prometheus-apple-%i/*'


# ── Adjust these to VPS ──────────────────────────────────────
User=adscrawler
WorkingDirectory=/home/adscrawler/adscrawler
ExecStart=/home/adscrawler/venvs/adscrawler-env/bin/python \
    -m dramatiq \
    adscrawler.dramatiq.app_stores.worker \
    -p 2 -t 4 \
    -Q store_crawls_apple_%i
# ──────────────────────────────────────────────────────────────────

Restart=always
RestartSec=10

# Hard and soft limits — tune if needed
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
