# adscrawler/dramatiq/crawl-store-apps-google@.service
# systemd template unit for a Google Play Dramatiq worker.
#
# Instantiate per country-priority group.  The group number (1 or 2) is
# passed via %i.
#
# Examples:
#
#   sudo systemctl enable --now crawl-store-apps-google@1
#   sudo systemctl enable --now crawl-store-apps-google@2
#
# Google Play scraping is sequential per chunk, so more processes (-p 4)
# provide parallelism.  The -t 2 gives two Dramatiq worker threads per
# process, each pulling its own chunk from Redis.
#
# The %i placeholder is expanded to the instance name (e.g. "1").

[Unit]
Description=Store crawler – Google group %i
After=network.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
Environment=dramatiq_prom_port=9191
Environment=dramatiq_prom_db=/tmp/dramatiq-prometheus-google-%i
Environment=prometheus_multiproc_dir=/tmp/dramatiq-prometheus-google-%i
ExecStartPre=/bin/sh -c 'mkdir -p /tmp/dramatiq-prometheus-google-%i && rm -rf /tmp/dramatiq-prometheus-google-%i/*'


# ── Adjust these to VPS ──────────────────────────────────────
User=adscrawler
WorkingDirectory=/home/adscrawler/adscrawler
ExecStart=/home/adscrawler/venvs/adscrawler-env/bin/python \
    -m dramatiq \
    adscrawler.dramatiq.app_stores.worker \
    -p 4 -t 2 \
    -Q store_crawls_google_%i
# ──────────────────────────────────────────────────────────────────

Restart=always
RestartSec=10

# Hard and soft limits — tune if needed
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
