# adscrawler/dramatiq/dramatiq-worker-google@.service
# systemd template unit for a Google Play Dramatiq worker.
#
# Instantiate per queue.  The queue name is passed via %i.
# Examples:
#
#   /etc/systemd/system/dramatiq-worker-google@store_crawls_google_1.service
#   → sudo systemctl enable --now dramatiq-worker-google@store_crawls_google_1
#
#   /etc/systemd/system/dramatiq-worker-google@store_crawls_google_2.service
#   → sudo systemctl enable --now dramatiq-worker-google@store_crawls_google_2
#
# Google Play scraping is sequential per chunk, so more processes (-p 4)
# provide parallelism.  The -t 2 gives two Dramatiq worker threads per
# process, each pulling its own chunk from Redis.
#
# The %i placeholder is expanded to the instance name (e.g.
# "store_crawls_google_1").

[Unit]
Description=Dramatiq Worker – %i
After=network.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1

# ── Adjust these to VPS ──────────────────────────────────────
User=adscrawler
WorkingDirectory=/home/adscrawler/adscrawler
ExecStart=/home/adscrawler/venvs/adscrawler-env/bin/python \
    -m dramatiq \
    adscrawler.dramatiq.app_stores.worker \
    -p 4 -t 2 \
    -Q %i
# ──────────────────────────────────────────────────────────────────

Restart=always
RestartSec=10

# Hard and soft limits — tune if needed
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
